# configurations/supervisord.conf
[supervisord]
nodaemon=true

[program:tyk]
command=/usr/bin/tyk --debug --conf=/etc/tyk/tyk.conf

[program:tyk-analytics]
command=/etc/tyk/tyk_dash/tyk-analytics
